# phoromatic.ini
[program:phoromatic]
nodaemon=true
command=phoronix-test-suite phoromatic.connect benchmark.paas.fasthosts.co.uk:80/CPEJEL
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
