# phoromatic.ini
[program:phoromatic]
nodaemon=true
command=/usr/bin/phoronix-test-suite start-phoromatic-server
user=root
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
